# milter-python2/etc/daemon.ini
[program:zimbra_mailinglists_milter]
command=/opt/zimbra_mailinglists_milter/zimbra_mailinglists_milter.py
process_name=zimbra_mailinglists_milter
priority=1
redirect_stderr=true
stdout_logfile=/var/log/zimbra_mailinglists_milter.log
